# etc/xdg-desktop-portal-jay.service
[Unit]
Description=Jay Portal

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.jay
ExecStart=/home/julian/bin/jay portal
